# hari-frontend-nginx-server/init/cron
#!/sbin/openrc-run
command=/usr/sbin/crond
command_args=""
command_user=root
command_background=yes
pidfile=/run/${RC_SVCNAME}.pid
