# hari-frontend-nginx-server/init/update-certificate
#!/sbin/openrc-run
command=/usr/local/bin/update-certificate
command_user=root